# datalab-ssh.timer
[Unit]
Description=Run datalab-ssh regularly
Requires=datalab-ssh.service

[Timer]
OnActiveSec=17m

[Install]
WantedBy=timers.target
